# backend/scripts/systemd/clipper-scraper.timer
[Unit]
Description=Run Clipper Clip Scraper Daily
Documentation=https://github.com/subculture-collective/clipper
Requires=clipper-scraper.service

[Timer]
# Run daily at 2:00 AM UTC
OnCalendar=*-*-* 02:00:00

# Alternative schedules (comment/uncomment as needed):
# Run every 6 hours:
# OnCalendar=*-*-* 00,06,12,18:00:00
#
# Run weekly on Sunday at 2 AM:
# OnCalendar=Sun *-*-* 02:00:00
#
# Run twice daily (2 AM and 2 PM):
# OnCalendar=*-*-* 02,14:00:00

# Persist timer across reboots
Persistent=true

# Run immediately if missed (e.g., system was down)
AccuracySec=1h

[Install]
WantedBy=timers.target
