# backend/scripts/systemd/clipper-scraper.service
[Unit]
Description=Clipper Targeted Clip Scraper
Documentation=https://github.com/subculture-collective/clipper
After=network-online.target postgresql.service redis.service
Wants=network-online.target
Requires=postgresql.service redis.service

[Service]
Type=oneshot
User=clipper
Group=clipper
WorkingDirectory=/opt/clipper/backend

# Environment setup
Environment="PATH=/usr/local/go/bin:/usr/bin:/bin"
EnvironmentFile=-/opt/clipper/backend/.env

# Main command
ExecStart=/opt/clipper/backend/bin/scrape_clips

# Resource limits
MemoryLimit=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clipper-scraper

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# Restart policy (for failures)
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target
